# deploy/hermes/hermes-mfp-bolus-sync.timer
[Unit]
Description=Run Hermes MyFitnessPal to Bolus AI nutrition sync every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
RandomizedDelaySec=60s

[Install]
WantedBy=timers.target
